# docker/overlay/etc/s6-overlay/s6-rc.d/init-tic-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Make folders
mkdir -p \
    /config/.tvh_iptv_config

# Set permissions
echo "Setting permissions"
lsiown -R abc:users /config/.tvh_iptv_config
