# docker/overlay/etc/s6-overlay/s6-rc.d/init-tic-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# make folders
mkdir -p \
    /config

# permissions
echo "Setting permissions"
lsiown -R abc:abc /config
